# infrastructure/wiki-dashboard.service
[Unit]
Description=WikiVelocity Analytics Dashboard
After=network.target

[Service]
ExecStart=/usr/bin/python3 /home/ubuntu/dashboard.py
Restart=always
RestartSec=5
User=ubuntu
Environment=AWS_DEFAULT_REGION=eu-west-1

[Install]
WantedBy=multi-user.target